# deploy/hypevision-api.service
[Unit]
Description=HypeVision API (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/hp-pan/backend
EnvironmentFile=-/root/hp-pan/backend/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/hp-pan/venv/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
