# deploy/hypevision-deploy-webhook.service
[Unit]
Description=HypeVision GitHub Deploy Webhook
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pnl
EnvironmentFile=-/home/pnl/deploy/deploy.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/pnl/deploy/webhook_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
